# roles/rke2_cluster/templates/vip-kube-apiserver-hc@.service.j2
{%- set hc_command = "/bin/true" -%}
{%- set vip_service_name = "vip-kube-apiserver@%i.service" -%}
[Unit]
Description=Health check for kube-apiserver VIP %i

After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Skipping health check if VIP service is inactive
ExecCondition=/bin/sh -c 'systemctl --quiet is-active {{ vip_service_name }}'

# Performing health check
# If health check passes, ensure VIP service is started
# If health check fails, ensure VIP service is stopped
ExecStart=/bin/sh -c '\
  if {{ hc_command }}; then \
    systemctl --quiet start {{ vip_service_name }}; \
  else \
    systemctl --quiet stop {{ vip_service_name }}; \
  fi'