# roles/rke2_cluster/templates/vip-kube-apiserver-hc@.timer.j2
[Unit]
Description=Periodic health check for kube-apiserver VIP %i

[Timer]
OnActiveSec=1s
OnUnitActiveSec=60s